# modules/apt/templates/apt-daily-upgrade.timer.erb
[Unit]
Description=Daily apt upgrade and clean activities
After=apt-daily.timer

[Timer]
<% if @autoupdate_times.length >= 2 -%>
OnCalendar=
OnCalendar=<%= @autoupdate_times[1] %>
RandomizedDelaySec=0
<% else -%>
OnCalendar=*-*-* 6:00
RandomizedDelaySec=60m
<% end -%>
Persistent=true

[Install]
WantedBy=timers.target
